# templates/profile/puppetboard/gunicorn.socket.epp
[Unit]
Description=puppetboard socket

[Socket]
ListenStream=<%= $bind %>
SocketUser=<%= $user %>

[Install]
WantedBy=sockets.target
